# templates/newlock.service.template
[Unit]
Description=SKABEN Lock service
After=network.target auditd.service

[Service]
ExecStart=${dirpath}/conf/run.sh
Restart=on-failure
Type=idle

[Install]
WantedBy=multi-user.target
